# tags/0.9.8/debian/prerm
#! /bin/sh
#
# This is the prerm script for the Debian GNU/Linux l2tp-ipsec-vpn-daemon package
#
# Written by Werner Jaeger
#
set -e

if [ "$1" = "remove" ] ; then
   # if daemon is running try to terminate it
   invoke-rc.d l2tp-ipsec-vpn-daemon stop || true

   # Remove runtime directory
   rm -rf /var/run/L2tpIPsecVpnControlDaemon
fi


# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
